# ops/launchd/com.windburn.multica-codex-cache-janitor.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.windburn.multica-codex-cache-janitor</string>

  <key>ProgramArguments</key>
  <array>
    <string>/Users/0xvox/Windburn/scripts/multica-codex-cache-janitor.sh</string>
    <string>--apply</string>
    <string>--min-age-minutes</string>
    <string>5</string>
  </array>

  <key>StartInterval</key>
  <integer>900</integer>

  <key>RunAtLoad</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/Users/0xvox/Library/Logs/windburn-multica-cache-janitor.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/0xvox/Library/Logs/windburn-multica-cache-janitor.err.log</string>
</dict>
</plist>
